# benchmark/GetApp/strictapps/htsn-0.1.1/doc/init.openrc
#!/sbin/runscript

#
# An example init file tailored for OpenRC on Gentoo.
#

description="Launch htsn to parse The Sports Network feed"
command="/usr/bin/htsn"
command_args="--daemonize"
pidfile="/run/htsn/htsn.pid"

depend() {
	use logger
	need net
}
